# pkg/DEBIAN/postinst
#!/bin/sh
set -e

case "$1" in
   configure)
      systemctl daemon-reload
      systemctl enable simple-bt.service
      systemctl start simple-bt.service
   ;;
esac

exit 0
